# supervisor/program/php-fpm.conf
[program:php-fpm]
autorestart=true
command=bash -c 'php-fpm7 -F'
priority=300